# etc-udev-rules.d/99-hid-apple.rules
ACTION=="add", SUBSYSTEM=="hid", DRIVER=="apple", RUN+="/bin/sh -c 'echo 2 > /sys/module/hid_apple/parameters/fnmode'"
